# src/ServiceMonitor/packaging/deb/postrm
#!/bin/bash
set -e
systemctl stop servicemonitor.service || true
systemctl disable servicemonitor.service || true
systemctl daemon-reload
